# contrib/crawl-scoreboard.service
[Unit]
Description=Crawl Scoreboard

[Service]
ExecStart=/app/scorebboard/contrib/run.sh

[Install]
WantedBy=multi-user.target
