# contrib/crawl-scoreboard.service
[Unit]
Description=Crawl Scoreboard

[Service]
ExecStart=/app/scoreboard/contrib/run.sh
Restart=on-success
RestartSec=5s

[Install]
WantedBy=multi-user.target
